# devstack/image/files/init-scripts/systemd/dpdkapp_pong.service
[Unit]
Description = dpdk app pong

[Service]
ExecStart = /usr/local/bin/l2_ping_pong -- --pong
Type = simple

